# Library/LaunchAgents/com.athal.ea-calendar.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.athal.ea-calendar</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>osascript {{ .chezmoi.homeDir }}/.config/opencode/skill/attention/ea/sync-calendars.applescript &amp;&amp; osascript {{ .chezmoi.homeDir }}/.config/opencode/skill/attention/ea/transition-blocks.applescript</string>
  </array>
  <key>StartInterval</key>
  <integer>1200</integer><!-- every 20 minutes -->
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/ea-calendar.log</string>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/ea-calendar.error.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>{{ .chezmoi.homeDir }}</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
