# Library/LaunchAgents/com.gh-pr-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gh-pr-poll</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-l</string>
        <string>-c</string>
        <string>source ~/.secret 2>/dev/null; ~/.local/bin/gh-pr-poll</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>StartInterval</key>
    <integer>120</integer>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/athal/.local/share/gh-pr-poll/poll.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/athal/.local/share/gh-pr-poll/poll.log</string>
</dict>
</plist>
