# Library/LaunchAgents/kb-enrich.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.{{ .chezmoi.username }}.kb-enrich</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>
REPO_DIR="{{ .chezmoi.homeDir }}/code/dotfiles"
SESSION_ID=$(/usr/bin/curl -s -X POST "http://localhost:4096/session?directory=$REPO_DIR" \
  -H "Content-Type: application/json" -d '{}' \
  | /usr/bin/python3 -c "import sys,json; print(json.load(sys.stdin).get('id',''))")
[ -z "$SESSION_ID" ] &amp;&amp; exit 1
/usr/bin/curl -s --max-time 10 -X POST "http://localhost:4096/session/$SESSION_ID/message?directory=$REPO_DIR" \
  -H "Content-Type: application/json" \
  -d '{"parts": [{"type": "text", "text": "Run the daily knowledge base enrichment for yesterday. Check activity across sources: opencode coding sessions, slack chat, zoom meeting transcripts, linear issues, and gh code reviews. Extract people facts, project updates, and decisions. Write a journal entry. Merge into knowledge-base profiles."}]}' || true
        </string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/kb-enrich.log</string>
</dict>
</plist>
